# deploy/Heartbeat-main.target
[Unit]
Description=Heartbeat Bootscript
Requires=Heartbeat-client.service Heartbeat-core.service

[Install]
WantedBy=multi-user.target
