# deploy/TSkr-main.target
[Unit]
Description=TSkr Bootscript
Requires=TSkr-client.service TSkr-core.service

[Install]
WantedBy=multi-user.target
